# supervisord.conf
[supervisord]
nodaemon=true
logfile=/tmp/supervisord.log
pidfile=/tmp/supervisord.pid

[program:xvfb]
command=Xvfb :1 -screen 0 %(ENV_SCREEN_WIDTH)sx%(ENV_SCREEN_HEIGHT)sx%(ENV_SCREEN_DEPTH)s -ac +extension GLX +render -noreset
autorestart=true
priority=10
startsecs=1
stdout_logfile=/tmp/xvfb.log
stderr_logfile=/tmp/xvfb.log

[program:x11vnc]
command=x11vnc -display :1 -rfbport %(ENV_VNC_PORT)s -nopw -forever -shared -noxdamage -noxfixes
autorestart=true
priority=20
startsecs=3
stdout_logfile=/tmp/x11vnc.log
stderr_logfile=/tmp/x11vnc.log

; noVNC servindo direto na porta 8080 com websockify embutido
; Isso elimina o serve.py como proxy e resolve o WebSocket
[program:novnc]
command=websockify --web /usr/share/novnc 8080 localhost:%(ENV_VNC_PORT)s
autorestart=true
priority=30
startsecs=4
stdout_logfile=/tmp/novnc.log
stderr_logfile=/tmp/novnc.log

[program:kelda]
command=/home/kelda/startup.sh
autorestart=false
priority=40
startsecs=10
stdout_logfile=/tmp/kelda.log
stderr_logfile=/tmp/kelda.log
